# deployment/teller-sync.timer
[Unit]
Description=Teller Transaction Sync Timer (Twice Daily)
Requires=teller-sync.service

[Timer]
# Run at 6 AM and 6 PM daily
OnCalendar=06:00
OnCalendar=18:00
Persistent=true

[Install]
WantedBy=timers.target